# recipes-graphics/xorg-xserver/xserver-xorg/X.service
[Unit]
Description=X Server

[Service]
Type=simple
ExecStartPre=-/bin/rm -rf /var/log/Xorg*
ExecStart=/usr/bin/X

[Install]
WantedBy=multi-user.target